# scripts/iot-server.service
[Unit]
Description=IOT server for room monitor

[Service]
ExecStart=/home/pi/main
TimeoutSec=30
Restart=on-failure
RestartSec=5
StartLimitInterval=350

[Install]
WantedBy=multi-user.target
